# data/upstart/multilatex.conf
#!upstart

description "multilatex node server"
author "Paul Nechifor <paul@nechifor.net>"

start on startup
stop on shutdown

# Run as an unprivileged user.
setuid multilatex
setgid multilatex

# Automatically respawn if it dies, but abandon if more than 5 times in 60 secs.
respawn
respawn limit 5 60

script
    exec /usr/bin/node /opt/multilatex/app/app.js >> /var/log/multilatex 2>&1
end script
